# scripts/com.jsperson.inbox-processor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jsperson.inbox-processor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jsperson/source/second_brain/scripts/process_inbox.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jsperson/Library/Mobile Documents/iCloud~md~obsidian/Documents/scott</string>

    <!-- Watch Inbox for changes (new files, modifications) -->
    <key>WatchPaths</key>
    <array>
        <string>/Users/jsperson/Library/Mobile Documents/iCloud~md~obsidian/Documents/scott/Second Brain/Inbox</string>
    </array>

    <!-- Prevent running more than once per 10 seconds (batches rapid changes) -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Fallback: also run periodically in case WatchPaths misses something -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- Also run immediately when loaded -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/jsperson/.imessage-capture/inbox-processor.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jsperson/.imessage-capture/inbox-processor-error.log</string>

    <!-- Environment for Claude -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/jsperson</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/jsperson/.npm-global/bin</string>
    </dict>
</dict>
</plist>
