# scripts/com.jsperson.weekly-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jsperson.weekly-review</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jsperson/.npm-global/bin/claude</string>
        <string>--print</string>
        <string>--dangerously-skip-permissions</string>
        <string>/weekly-review</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jsperson/source/second_brain</string>

    <!-- Run at 16:00 on weekday 0 (0=Sunday) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>16</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/jsperson/.imessage-capture/weekly-review.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jsperson/.imessage-capture/weekly-review-error.log</string>

    <!-- Environment for Claude -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/jsperson</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/jsperson/.npm-global/bin</string>
    </dict>
</dict>
</plist>
